# deploy/api-private-router.service
[Unit]
Description=api-private-router - Application Runtime
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=api-private-router
Group=api-private-router
WorkingDirectory=/opt/api-private-router
ExecStart=/usr/bin/java -jar /opt/api-private-router/api-private-router.jar
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=api-private-router

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/api-private-router /etc/api-private-router

# Environment - Load from .env file (create /etc/api-private-router/env)
EnvironmentFile=-/etc/api-private-router/env

# Server configuration (defaults, can be overridden in env file)
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080
Environment=DATA_DIR=/opt/api-private-router
Environment=APP_BUILD_TYPE=release
Environment=DATABASE_HOST=127.0.0.1
Environment=DATABASE_PORT=5432
Environment=REDIS_HOST=127.0.0.1
Environment=REDIS_PORT=6379

[Install]
WantedBy=multi-user.target
